# archiso-autobuild/etc/systemd/system/buildarchiso.service
[Unit]
Description=Build Archiso

[Service]
User=root
Type=oneshot
ExecStart=+/root/buildarchiso.sh
