# archiso-autobuild/etc/systemd/system/buildarchiso.timer
[Unit]
Description=Build Archiso daily

[Timer]
OnCalendar=0:00
Persistent=true

[Install]
WantedBy=timers.target
